# startup.service
[Unit]
Description=nginx reload

[Service]
ExecStart=/script.sh

[Install]
WantedBy=multi-user.target
